# supervisor/usb-backup.conf
[program:usb-backup]
directory = /usr/share/usb-backup
environment =
command = python usb-backup.py
autostart = true
